# static/etc/services.d/couchpotato/run
#!/usr/bin/with-contenv bash

umask 0002

chown -R bytesized:bytesized /app /movies /config

exec s6-setuidgid bytesized python /app/couchpotato/CouchPotato.py --data_dir /config
